# volcengine-assist-client/static/etc/systemd/system/install-volc-assist-client.service
[Unit]
Description=Try to install volcengine assist client

[Service]
Type=oneshot
ExecStart=/opt/assist-client/check-install.sh

[Install]
WantedBy=multi-user.target
