# aliyun-assist-client/static/etc/systemd/system/install-assist-client.service
[Unit]
Description=Try to install assist client

[Service]
Type=oneshot
ExecStart=/opt/assist-client/check-install.sh

[Install]
WantedBy=multi-user.target
